# cron/assets/crontab
* * * * * /cron-script.sh 2>>/log >> /log
